# certidude/builder/overlay/etc/uci-defaults/90-certidude-sysupgrade
#!/bin/sh
echo /etc/certidude/authority/ >> /etc/sysupgrade.conf

