# certidude/templates/server/ldap-kinit.timer
[Timer]
OnCalendar=00/8:30
Persistent=true
