# certidude/templates/server/ldap-kinit.timer
[Timer]
Persistent=true
OnBootSec=3sec
OnUnitActiveSec=8h

[Install]
WantedBy=multi-user.target
