# certidude/templates/client/certidude.service
[Unit]
Description=Renew certificates and update revocation lists

[Service]
Type=simple
ExecStart={{ sys.argv[0] }} request
