# certidude/templates/client/certidude.timer
[Unit]
Description=Run certidude service weekly

[Timer]
OnCalendar=weekly
Persistent=true
Unit=certidude.service

[Install]
WantedBy=timers.target

