# certidude/templates/snippets/nginx-ocsp-cache.timer
[Timer]
OnCalendar=*:0/15
Persistent=true
